# docker/wolfstub/docker/supervisor/supervisor-lemp.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize
priority=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=2
